# init/cybermon-cassandra.service.in

[Unit]
Description=cybermon Cassandra service
Documentation=https://cybermaggedon.github.io/cyberprobe-docs/

[Service]
PIDFile=/var/run/cybermon-cassandra.pid
ExecStart=@bindir@/cybermon-cassandra

[Install]
WantedBy=multi-user.target

