# init/cyberprobe.service.in

[Unit]
Description=cyberprobe service
Documentation=http://cyberprobe.sourceforge.net/cyberprobe.html
Requires=cybermon.service

[Service]
PIDFile=/var/run/cyberprobe.pid
ExecStart=@bindir@/cyberprobe @sysconfdir@/cyberprobe.cfg

[Install]
WantedBy=multi-user.target

