# init/cybermon.service.in

[Unit]
Description=cybermon service
Documentation=http://cyberprobe.sourceforge.net/cyberprobe.html

[Service]
PIDFile=/var/run/cybermon.pid
ExecStart=@bindir@/cybermon -p 9000 -c @sysconfdir@/@PACKAGE@/cybermon.lua

[Install]
WantedBy=multi-user.target

