# init/cyberprobe.service.in

[Unit]
Description=cyberprobe service
Documentation=https://cybermaggedon.github.io/cyberprobe-docs/

[Service]
PIDFile=/var/run/cyberprobe.pid
ExecStart=@bindir@/cyberprobe @sysconfdir@/cyberprobe.cfg

[Install]
WantedBy=multi-user.target

