# init/cybermon-elasticsearch.service.in

[Unit]
Description=cybermon ElasticSearch service
Documentation=http://cyberprobe.trustnetworks.com/cyberprobe.html

[Service]
PIDFile=/var/run/cybermon-bigquery.pid
ExecStart=@bindir@/cybermon-elasticsearch

[Install]
WantedBy=multi-user.target

