# init/evs-gaffer.service.in

[Unit]
Description=cyberprobe Gaffer loader service
Documentation=https://cybermaggedon.github.io/cyberprobe-docs/

[Service]
PIDFile=/var/run/cybermon-gaffer.pid
ExecStart=@bindir@/cybermon-gaffer

[Install]
WantedBy=multi-user.target

