# init/evs-elasticsearch.service.in

[Unit]
Description=cyberprobe ElasticSearch loader service
Documentation=https://cybermaggedon.github.io/cyberprobe-docs/

[Service]
PIDFile=/var/run/cybermon-bigquery.pid
ExecStart=@bindir@/cybermon-elasticsearch

[Install]
WantedBy=multi-user.target

