# init/evs-geoip.service.in

[Unit]
Description=cyberprobe IP-to-geo service
Documentation=https://cybermaggedon.github.io/cyberprobe-docs/

[Service]
PIDFile=/var/run/cybermon-geoip.pid
ExecStart=@bindir@/cybermon-geoip

[Install]
WantedBy=multi-user.target

