# cyberprobe-group.service

[Unit]
Description=cyberprobe group

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true

